# defaults/scala.service
[Unit]
Description=Scala Daemon
After=network.target

[Service]
Type=forking
GuessMainPID=no
ExecStart=/home/pool/scala/build/release/bin/scalad \
--rpc-bind-ip 127.0.0.1 \
--rpc-bind-port 20189 \
--p2p-bind-ip 0.0.0.0 \
--p2p-bind-port 20188 \
--rpc-bind-ip 127.0.0.1 \
--confirm-external-bind --detach

Restart=always
User=pool

[Install]
WantedBy=multi-user.target
